# ibenc.timer
[Unit]
Description=Run iBench every hour
Requires=ibenc.service

[Timer]
# Run every hour at the top of the hour
OnCalendar=hourly

# Or uncomment below to run at specific times (e.g., 12:00, 18:00, 00:00, 06:00)
# OnCalendar=*-*-* 00:00:00,06:00:00,12:00:00,18:00:00

# Persist timer state across reboots
Persistent=true

# Run immediately if the system was powered off during scheduled time
AccuracySec=1min

[Install]
WantedBy=timers.target
