# ibenc.service
[Unit]
Description=iBench
Documentation=https://github.com/yourusername/ibenc
After=network-online.target
Wants=network-online.target

[Service]
Type=oneshot
ExecStart=/usr/local/bin/ibenc -config /etc/ibenc/ibenc.yaml

# Environment and logging
User=root
StandardOutput=journal
StandardError=journal
SyslogIdentifier=ibenc

# Timeout for test execution (tests can take 30-60 seconds)
TimeoutStartSec=300

[Install]
WantedBy=multi-user.target
